# etc/cron.daily/syncVideos
#!/bin/bash

videosFolder=/home/pi/minidlna/VideosRoot/Videos

/home/pi//minidlna/scripts/syncVideos.sh $videosFolder
